# infra/ansible/roles/hermes/templates/hermes-dashboard.service.j2
[Unit]
Description=Hermes Agent web dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ hermes_runtime_user | default('anvil') }}
Group={{ hermes_runtime_user | default('anvil') }}
EnvironmentFile=/etc/hermes-dashboard.env
WorkingDirectory=/home/{{ hermes_runtime_user | default('anvil') }}
{% if hermes_managed_runtime_supported | default(false) | bool %}
ExecStartPre=/usr/local/libexec/hermes-dashboard-preflight
{% endif %}
ExecStart=/usr/local/bin/hermes dashboard --host {{ hermes_dashboard_host | default('127.0.0.1') }} --port {{ hermes_dashboard_port | default(9119) }} --no-open --insecure
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
