# ansible/roles/caddy_proxy/templates/ssh-proxy.socket.j2
[Unit]
Description={{ item.name | default('SSH') }} proxy socket

[Socket]
ListenStream={{ item.listen_port }}
NoDelay=true

[Install]
WantedBy=sockets.target
